# com.hyds.realtime.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hyds.realtime</string>

    <key>WorkingDirectory</key>
    <string>__HYDS_PATH__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__HYDS_PATH__/venv/bin/python</string>
        <string>execution/realtime_check.py</string>
    </array>

    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>__HYDS_PATH__/.tmp/realtime.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HYDS_PATH__/.tmp/realtime.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
